# install-files/usr/lib/NetworkManager/dispatcher.d/01-lliurex-autotags-network
#!/usr/bin/env bash

case "$2" in
    up)
        /usr/sbin/lliurex-autotags-network
    ;;

    down)
        /usr/sbin/lliurex-autotags-network
    ;;

    *)
    ;;
esac
