# install-files/usr/lib/systemd/system/lliurex-system-tags.service
[Unit]
Description=System tags for autoupgrades

[Service]
Type=oneshot
ExecStart=/usr/sbin/lliurex-system-tags
ExecStop=/bin/true

[Install]
WantedBy=default.target
