# install-files/lib/systemd/system/lliurex-autotags-network.service
[Unit]
Description=Network tags for autoupgrades
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/lliurex-autotags-network
ExecStop=/bin/true

[Install]
WantedBy=default.target
